# deb/bambooslacking_1.0-1/DEBIAN/postrm
#!/bin/bash

set -e

systemctl disable bambooslacking
systemctl daemon-reload